# .config/autostart/leftwm-settings.desktop
[Desktop Entry]
Comment=Load leftwm settings
Exec=leftwm-settings
Name=leftwm-settings
Terminal=true
Type=Application
